# .carnotzet/assets/rootfs/root/etc/sv/wifi/finish
#!/bin/sh

modprobe -r brcmfmac

sleep 5