# .carnotzet/assets/rootfs/root/etc/sv/wifi/finish
#!/bin/sh

modprobe -qr brcmfmac

sleep 5