# whispi+auto-shutdown-rpi64/assets/sv/log/run
#!/bin/sh
exec svlogd -tt /var/log/tunnel-autossh