# carnotzet+rootfs-rpi4b/root/etc/sv/network/run
#!/bin/sh
exec 2>&1

echo "Start network..."

# Charger la configuration
if [ -f /boot/carnotzet.cfg ]; then
  . /boot/carnotzet.cfg
else
  echo "Error: /boot/carnotzet.cfg not found"
  exit 1
fi

# Vérifier si interface réseau configurée
if [ -z "$network_wire" ]; then
  echo "Network interfaces disabled"
  exit 0
fi

echo "Configure interface $network_wire..."

if [ "$network_wire_dhcp" = "off" ]; then
  echo "Configure static IP: $network_wire_ip/$network_wire_mask"
  ip addr add "$network_wire_ip/$network_wire_mask" dev "$network_wire"
  ip link set dev "$network_wire" up
  ip route add default via "$network_wire_gw" dev "$network_wire"
  
  # Supervision
  while true; do
    sleep 300
    echo "Interface $network_wire enabled - $(date)"
  done
fi

if [ "$network_wire_dhcp" = "on" ]; then
  echo "Configure DHCP for $network_wire"
  ip link set dev "$network_wire" up
  
  # Run in foreground for runit
  exec udhcpc -f -i "$network_wire"
fi

echo "Error: unknown network settings"
exit 1